# images/drbd-sysext/sysext/sysext-drbd.service
# Fires once the drbd sysext is merged — a delivery-mechanism smoke marker.
# Proves the sysext was fetched by nodelet and merged by systemd-sysext (the
# marker file only exists after merge). The real unit (loading drbd kernel
# modules) replaces this once the compiled drbd.ko payload lands.
[Unit]
Description=DRBD sysext delivery marker
After=systemd-sysext.service
ConditionPathExists=/usr/share/drbd-sysext/marker

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c 'echo "drbd sysext merged: $(cat /usr/share/drbd-sysext/marker)"'

[Install]
WantedBy=multi-user.target
