# images/drbd-sysext/sysext/sysext-drbd.service
# Indexes the DRBD kernel modules shipped by this sysext (drbd.ko /
# drbd_transport_tcp.ko under /usr/lib/modules/<ver>/updates/) so modprobe can
# resolve them after the sysext is merged. nodelet performs the actual load of
# `drbd` / `drbd_transport_tcp` from NodeExtensionRequest.kernelModules.
[Unit]
Description=Index DRBD kernel modules from the drbd sysext
After=systemd-sysext.service
ConditionPathExists=/usr/lib/modules/%v/updates/drbd.ko

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/depmod %v

[Install]
WantedBy=multi-user.target
